# etc/services.d/init-waitd/run
#!/usr/bin/execlineb -S0
# script used to notify scripts that the container init scripts finished
# this service is a run once:

# This service wishes to stay down after this run, thank you
foreground { touch /var/run/s6/services/init-waitd/down }
foreground { s6-svc -O /var/run/s6/services/init-waitd }

importas -D "/var/run/s6/.init-finished" S6_INIT_FINISHED_FILE S6_INIT_FINISHED_FILE
importas -D "/var/run/s6/.init-wait" S6_INIT_FINISHED_EVENT S6_INIT_FINISHED_EVENT

# solve the chicken-egg problem using a separate file that will be checked
# by the subscriber to see if the event has already been notified
if { touch ${S6_INIT_FINISHED_FILE} }

# Create the event dir (if doesn't exist) and send the notification event
foreground { s6-mkfifodir "${S6_INIT_FINISHED_EVENT}" }
foreground { s6-ftrig-notify "${S6_INIT_FINISHED_EVENT}" x }

