# etc/s6-overlay/s6-rc.d/init-timezone/up
# Automatically configures the timezone if the TZ variable is set
with-contenv

importas -D "" TZ TZ
if -t { eltest -n "${TZ}" }

if { ln -snf "/usr/share/zoneinfo/${TZ}" /etc/localtime }

redirfd -w 1 /etc/timezone s6-echo "${TZ}"

